# configs/systemd/.config/systemd/user/uair.service
[Unit]
Description=Uair

[Service]
Type=simple
ExecStart=/usr/local/bin/uair

[Install]
WantedBy=default.target
